# cmd/namestr.service
[Unit]
Description=namestr

[Service]
Type=simple
User=me
ExecStart=/home/me/.local/bin/namestrt

[Install]
WantedBy=multi-user.target